# trance-daemon/debian/postrm
#!/bin/sh
# postrm script for idlescreen-daemon
set -e

if [ "$1" = "purge" ]; then
    # Purge the user service enable state globally using standard deb-systemd-helper
    if command -v deb-systemd-helper >/dev/null; then
        deb-systemd-helper --user purge idlescreen-daemon.service || true
    fi
fi

exit 0
